# daemon/claude-relay.service
[Unit]
Description=Discord + Telegram Relay (Claude Code)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discord-telegram-relay
ExecStart=/usr/bin/bun run /home/ubuntu/discord-telegram-relay/src/relay.js
Restart=always
RestartSec=10

# Environment variables (set these separately or use a .env file)
Environment=TELEGRAM_BOT_TOKEN=your_telegram_token
Environment=TELEGRAM_USER_ID=your_telegram_id
Environment=DISCORD_BOT_TOKEN=your_discord_token
Environment=DISCORD_CHANNEL_IDS=your_channel_ids
Environment=MCP_URL=http://localhost:3001/sse
Environment=MCP_API_KEY=your_mcp_api_key
Environment=CLAUDE_PATH=claude
Environment=USER_NAME=Diego
Environment=USER_TIMEZONE=America/Lima

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/.discord-telegram-relay

[Install]
WantedBy=multi-user.target
